# mudpi.service
[Unit]
Description=Mudpi server
after=network.target

[Service]
User=mud
ExecStart=/home/mudpi/mudpi/start.sh
WorkingDirectory=/home/mudpi/mudpi
Restart=on-failure

[Install]
WantedBy=multi-user.target
